# app/modules/queclink/gl300/gl300_proxy.conf
# # GL300 proxy service -- save this in /etc/init/gl300_proxy.conf
# # service start gl300_proxy

description     "GL300 Proxy Server"
author          "GPS.Tools"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
exec /usr/bin/ruby /home/raco1/domains/ror.raco1.websitesonwheels.net/public_html/app/modules/queclink/gl300/proxy.rb
